# templates/etc/systemd/system/pgbackrest-backup.timer.j2
[Unit]
Description=PgBackRest Backup Job

[Timer]
OnCalendar={{ pgbackrest_systemd_timer_oncalendar }}
Persistent=true
Unit={{ pgbackrest_systemd_timer_unit }}

[Install]
WantedBy=timers.target
