# templates/pgbackrest-backup.timer.j2
[Unit]
Description=PgBackRest Backup Job

[Timer]
OnCalendar={{ pgbackrest_systemd_timer_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
